# configs/rootfs.template/etc/rc.d/rc.local
#!/bin/sh

echo "start http_server"
cd /opt/bin && nohup /opt/bin/http_server &
#cd /opt/bin/ && /opt/bin/http_server > http_server.log 2>&1 &
echo "finish rc.local"
